# e20-kit1-demo/e20-sys/e20-kit1-demo.conf
# Synapse example - start as a service

description	"Start e20-kit1-demo web app"

start on runlevel [2]
stop on runlevel [!2]

exec python /home/snap/e20-kit1-demo/web_app/app_server.py

